# crontab
0 * * * * find /app/temp_results -type f -mmin +1440 -delete
